# devenv/docker/blocks/graphite1/conf/etc/service/nginx/run
#!/bin/bash

mkdir -p /var/log/nginx
exec /usr/sbin/nginx -c /etc/nginx/nginx.conf
