# devenv/docker/blocks/graphite1/conf/etc/service/statsd/run
#!/bin/bash

exec /usr/bin/nodejs /opt/statsd/stats.js /opt/statsd/config.js >> /var/log/statsd.log 2>&1

